# meta-hdmi-udev/recipes-core/udev/udev-hdmi-weston/99-weston-hdmi.rules
#!/bin/bash

SUBSYSTEM=="drm", ACTION=="change", ENV{HOTPLUG}=="1", RUN+="/etc/udev/scripts/weston_restart.sh"
